# config/rubies.service
[Unit]
Description=rubies
After=syslog.target network.target
Requires=redis-server.service
After=redis-server.service

[Service]
Type=simple
User=epistrephein
SyslogIdentifier=rubies
WorkingDirectory=/var/www/rubies.io
Environment="RACK_ENV=production"
Environment="REDIS_URL=redis://localhost:6379/0"
ExecStart=/bin/bash -l -c 'bundle exec puma -C config/puma.rb >> log/puma.log 2>&1'
Restart=always
TimeoutSec=10

[Install]
WantedBy=multi-user.target
